# foo-src/foo-daemon.service.in
[Unit]
Description=Sample daemon that demonstrates SCLs

[Service]
Type=simple
ExecStart=@bindir@/foo-daemon

[Install]
WantedBy=multi-user.target

